# airootfs/etc/profile.d/checkmate.sh
#!/bin/bash


echo "Apple loves chess"
echo ""

while true
do
	/usr/local/bin/checkra1n
done
